# wifi.target
[Unit]
Description=Robot wifi target
Requires=multi-user.target
